# com.magiccardscanner.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.magiccardscanner.server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/billjackson/Downloads/working/experiments/projects/magic-card-scanner-production && source venv/bin/activate && python -m uvicorn backend.app:app --host 0.0.0.0 --port 8000</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/billjackson/Downloads/working/experiments/projects/magic-card-scanner-production/logs/server.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/billjackson/Downloads/working/experiments/projects/magic-card-scanner-production/logs/server.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/billjackson/Downloads/working/experiments/projects/magic-card-scanner-production</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/billjackson/Downloads/working/experiments/projects/magic-card-scanner-production</string>
        <key>ENV_MODE</key>
        <string>production</string>
        <key>PORT</key>
        <string>8000</string>
    </dict>
</dict>
</plist> 